# com.emergent.kg2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.emergent.kg2</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/rocky/emergent/evolve-auto-kg2.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>WorkingDirectory</key>
    <string>/Users/rocky/emergent</string>
    <key>StandardOutPath</key>
    <string>/Users/rocky/emergent/kg2/logs/launchagent-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rocky/emergent/kg2/logs/launchagent-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/rocky</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
